# .containers/development/rootfs/etc/profile.d/umask.sh
# Set umask to allow groups write-permission.
umask 002
